# systemd/rinse-forward.service
[Unit]
Requires=rinse.service
After=rinse.service
Requires=rinse-forward.socket
After=rinse-forward.socket

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd 127.0.0.1:8443
